# projekteserver/backup/jira-backup.timer
[Unit]
Description=Backup JIRA Database and data

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
wantedBy=timers.target
