# projekteserver/backup/jira-backup.service
[Unit]
Description=Backup JIRA Database and data

[Service]
Type=simple
ExecStart=/opt/backups/backup.sh
